# core/s3c-init/inittab
::sysinit:/etc/rc.sysinit
::respawn:/sbin/login
