# core/s3c-init/inittab
::sysinit:/etc/rc.sysinit
::respawn:/usr/bin/monit -Ic /etc/monitrc
::respawn:/sbin/login
::ctrlaltdel:/usr/bin/monit stop all
::ctrlaltdel:/bin/umount -a -r
::ctrlaltdel:/bin/swapoff
