# core/s3c-init/inittab
::sysinit:/etc/rc.sysinit
::respawn:/usr/bin/monit -Ic /etc/monitrc
::respawn:/sbin/login
