# core/s3c-init/rc.local
# Put your own things to run on startup here
